# php-5/docker-utils/supervisord/php5-fpm.conf
[program:php5-fpm]
command=/usr/sbin/php5-fpm -F
stdout_events_enabled=true
stderr_events_enabled=true
autostart=true
autorestart=true

