# php-7/docker-utils/supervisord/apache2.conf
[program:apache2]
command=/usr/sbin/apache2ctl -D FOREGROUND
stdout_events_enabled=true
stderr_events_enabled=true
autostart=true
autorestart=true
